# mcp/systemd/em-unetpi-mcp.service.sample
[Unit]
Description=EM-uNetPi MCP Server
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/EM-uNetPi/mcp
ExecStart=/home/pi/.local/bin/uv run server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target